# files/system/etc/profile.d/mimalloc.sh
#!/usr/bin/sh
export LD_PRELOAD=libmimalloc.so